# com.example.daily-meeting-minutes.plist
<?xml version="1.0" encoding="UTF-8"?>
<!-- 参考/手动兜底模板。正常情况 ./init.sh 已根据你选择的频率自动生成并加载本文件到 -->
<!-- ~/Library/LaunchAgents/,无需手改。只有想手动接管时才用本模板:把 __SHARE_DIR__ -->
<!-- 换成本包绝对路径、__BASE_DIR__ 换成本地归档目录(默认 ~/会议纪要)、 -->
<!-- __PATH__ 换成你交互 shell 里 `echo $PATH` 的结果(关键!launchd 默认 PATH 极简, -->
<!-- 不补这个会找不到 lark-cli/sc/node/python3 而静默失败)。 -->
<!-- 频率二选一:每天定点用下面的 StartCalendarInterval;每隔 N 秒跑一次则改用 -->
<!-- <key>StartInterval</key><integer>N秒</integer>(如每 4 小时 = 14400)。Linux 用 crontab,见 README。 -->
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.example.daily-meeting-minutes</string>
    <key>ProgramArguments</key>
    <array>
        <string>/bin/zsh</string>
        <string>__SHARE_DIR__/run.sh</string>
    </array>
    <key>EnvironmentVariables</key>
    <dict>
        <key>PATH</key>
        <string>__PATH__</string>
    </dict>
    <key>WorkingDirectory</key>
    <string>__SHARE_DIR__</string>
    <key>StartCalendarInterval</key>
    <dict>
        <key>Hour</key><integer>18</integer>
        <key>Minute</key><integer>47</integer>
    </dict>
    <key>StandardOutPath</key>
    <string>__BASE_DIR__/_logs/launchd.out.log</string>
    <key>StandardErrorPath</key>
    <string>__BASE_DIR__/_logs/launchd.err.log</string>
    <key>RunAtLoad</key><false/>
</dict>
</plist>
